# fpf.timer
[Unit]
Description=Periodically sync files database

[Timer]
OnUnitInactiveSec=24hours

[Install]
WantedBy=multi-user.target
